# services.d/gluster/run
#!/usr/bin/with-contenv sh

chown -R gluster:gluster /var/log/glusterfs

exec /usr/sbin/glusterd -N
